# deploy/feishu-claude.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zara.feishu-claude</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/amber/Developer/feishu-claudecode-bridge/.venv/bin/python3</string>
        <string>-u</string>
        <string>/Users/amber/Developer/feishu-claudecode-bridge/main.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/amber/Developer/feishu-claudecode-bridge</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/amber/.local/share/fnm/node-versions/v24.15.0/installation/bin:/Users/amber/.local/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/amber</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/amber/Developer/feishu-claudecode-bridge/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/amber/Developer/feishu-claudecode-bridge/stderr.log</string>
</dict>
</plist>
